# dist/xyz.openbmc_project.Logging.AuditLog.service
[Unit]
Description=Audit Log Services

[Service]
ExecStart=/usr/bin/phosphor-auditlog
Restart=always
Type=dbus
BusName=xyz.openbmc_project.Logging.AuditLog

[Install]
WantedBy=multi-user.target
